# deploy/eve-t1-calc.service
[Unit]
Description=EVE T1 manufacturing profit calculator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eve
WorkingDirectory=/opt/eve-t1-calc
ExecStart=/opt/eve-t1-calc/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
